# yunxingcloud/yunxingcloud.service
[Unit]
Description=yunxingcloud Application
After=network.target mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/yunxingcloud/app
ExecStart=/usr/bin/java \
    -Xms512m -Xmx1024m \
    -XX:+UseG1GC \
    -Duser.timezone=Asia/Shanghai \
    -jar /opt/yunxingcloud/app/yunxingcloud-core/target/yunxingcloud-core-0.0.1-SNAPSHOT.jar \
    --server.port=8080 \
    --spring.profiles.active=prod
ExecStop=/bin/kill -15 $MAINPID
Restart=always
RestartSec=10
StandardOutput=append:/opt/yunxingcloud/logs/stdout.log
StandardError=append:/opt/yunxingcloud/logs/stderr.log

# 环境变量（通过 EnvironmentFile 注入）
EnvironmentFile=/opt/yunxingcloud/app/env.conf

[Install]
WantedBy=multi-user.target
