# yunxingcloud/yunxingcloud.service
[Unit]
Description=yunxingcloud Microservice (%i)
After=network.target mysql.service redis.service

[Service]
Type=simple
User=yunxingcloud
WorkingDirectory=/opt/yunxingcloud/app
ExecStart=/usr/bin/java \
    -Xms256m -Xmx512m \
    -XX:+UseG1GC -XX:MaxGCPauseMillis=200 \
    -XX:+HeapDumpOnOutOfMemoryError -XX:HeapDumpPath=/opt/yunxingcloud/logs/ \
    -Duser.timezone=Asia/Shanghai \
    -jar /opt/yunxingcloud/app/yunxingcloud-%i.jar \
    --server.port=%p \
    --spring.profiles.active=prod
ExecStop=/bin/kill -15 $MAINPID
Restart=on-failure
RestartSec=10
TimeoutStopSec=30

EnvironmentFile=/opt/yunxingcloud/app/env.conf

StandardOutput=append:/opt/yunxingcloud/logs/%i-stdout.log
StandardError=append:/opt/yunxingcloud/logs/%i-stderr.log

[Install]
WantedBy=multi-user.target